# AirPollutionControl/raspberry/dust_sensor/code_examples/github_sorz_codeexamples_dht22_gp2y1010/air-graph@.service
[Unit]
Description=Execute a script to (re)draw RRD graphs once.

[Service]
Type=simple
User=xierch
Group=xierch
ExecStart=/home/xierch/arduino/graph.sh %i
