# AirPollutionControl/raspberry/dust_sensor/code_examples/github_sorz_codeexamples_dht22_gp2y1010/air-graph-1w.timer
[Unit]
Description=Redraw air graphs.

[Timer]
OnBootSec=3min
OnUnitActiveSec=3h
Unit=air-graph@1w.service

[Install]
WantedBy=multi-user.target