# AirPollutionControl/raspberry/dust_sensor/code_examples/github_sorz_codeexamples_dht22_gp2y1010/rrdupdate-air.service
[Unit]
Description=Run rrdupdate for air data collecting.

[Service]
Type=simple
ExecStart=/home/xierch/arduino/arduino-uart.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
